# rmtfs-dir.service.in
[Unit]
Description=Qualcomm remotefs service
Before=NetworkManager.service

[Service]
ExecStart=RMTFS_PATH/rmtfs -s -o RMTFS_EFS_PATH
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
